# unms/root/etc/services.d/postgres/run
#!/usr/bin/with-contenv sh
s6-setuidgid postgres postgres -c listen_addresses='127.0.0.1'

